# siga-client-conf/gcs/newfiles_skel/etc/profile.d/move-ull-home.sh
if echo $HOME | grep -q ull.es; then
	export HOME="$HOME/home"
fi
