# siga-client-conf/gcs/newfiles_skel/etc/network/if-up.d/enable-wol
#!/bin/sh

# quit if we're called for the loopback
if [ "$IFACE" = lo ]; then
	exit 0
fi

/usr/sbin/ethtool -s $IFACE wol g
